# deployment/systemd/{{PROJECT_NAME}}.service
[Unit]
Description=Simple SNMP Daemon - A lightweight and secure SNMP monitoring daemon
After=network.target

[Service]
Type=simple
User=simple-snmpd
Group=simple-snmpd
ExecStart=/usr/local/bin/simple-snmpd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-snmpd

[Install]
WantedBy=multi-user.target
